# packaging/systemd/postgresql@apps.service.d/30-postgresql-setup.conf
[Service]
Environment=PGDATA=/var/lib/pgsql/apps/data
