# packaging/systemd/postgresql@apps.service.d/40-pgdg18-binary.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/pgsql-18/bin/postgresql-18-check-db-dir ${PGDATA}
ExecStart=
ExecStart=/usr/pgsql-18/bin/postgres -D ${PGDATA}
