# packaging/systemd/simplehost-control.service
[Unit]
Description=SimpleHost Control Plane
After=network-online.target postgresql@shp.service
Wants=network-online.target
Requires=postgresql@shp.service

[Service]
Type=simple
User=simplehost
Group=simplehost
WorkingDirectory=/opt/simplehostman/release/current
EnvironmentFile=/etc/simplehost/control.env
ExecStart=/usr/bin/node /opt/simplehostman/release/current/apps/control/dist/index.js
Restart=always
RestartSec=5
StateDirectory=simplehost
LogsDirectory=simplehost

[Install]
WantedBy=multi-user.target
