# packaging/systemd/simplehost-worker.service
[Unit]
Description=SimpleHost Worker
After=network-online.target postgresql@control.service
Wants=network-online.target
Requires=postgresql@control.service

[Service]
Type=simple
User=simplehost
Group=simplehost
WorkingDirectory=/opt/simplehostman/release/current
EnvironmentFile=/etc/simplehost/worker.env
ExecStart=/usr/bin/node /opt/simplehostman/release/current/apps/worker/dist/index.js
Restart=always
RestartSec=5
StateDirectory=simplehost
LogsDirectory=simplehost

[Install]
WantedBy=multi-user.target
