# packaging/systemd/postgresql@control.service.d/30-postgresql-setup.conf
[Service]
Environment=PGDATA=/var/lib/pgsql/control/data
