# platform/containers/quadlet/worker-template.container
[Unit]
Description=Worker container for __APP_SLUG__
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=worker-__APP_SLUG__
Image=__APP_IMAGE__
EnvironmentFile=/etc/containers/systemd/env/app-__APP_SLUG__.env
Volume=/srv/containers/apps/__APP_SLUG__/app:/var/www/html:Z
Volume=/srv/containers/apps/__APP_SLUG__/uploads:/var/www/html/public/uploads:Z
Exec=__WORKER_COMMAND__

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
